# templates/deluge/deluge.conf
start on (filesystem and networking) or runlevel [2345]
stop on runlevel [016]

env uid=USER
env gid=GROUP
env umask=022

exec start-stop-daemon -S -c $uid:$gid -k $umask -x /usr/bin/deluged -- -d
