# templates/sonarr/sonarr.conf
start on (filesystem and networking) or runlevel [2345]
stop on runlevel [016]

env uid=USER
env gid=GROUP

exec mono /opt/NzbDrone/NzbDrone.exe
