# templates/jackett/jackett.conf
start on runlevel [2345]
stop on runlevel [016]

env uid=USER
env gid=GROUP
env umask=022
env DIR=/opt/jackett

script
    chdir $DIR
    exec mono $DIR/Jackett.exe
end script
